# scripts/services/openvpn.conf
[program:openvpn]
command=/sbin/tini -- /usr/bin/openvpn.sh
